# root/rc.local
#!/bin/sh 
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
gsmctl -S -s "5550506471 Prueba de arranque"
/usr/bin/lua /root/Sistema.lua &
cp /root/Config.bak /root/Config.txt
exit 0
